# Terraform-Ansible-automation/ansible_diplom/app_basic_setup/files/media_backup.timer
[Unit]
Description=Media Backup timer every  day

[Timer]
OnCalendar=00:00:00
Unit=media_backup.service

[Install]
WantedBy=timers.target
