# Terraform-Ansible-automation/ansible_diplom/app_basic_setup/files/pgsql_backup.timer
[Unit]
Description=Whole pgsql backup timer

[Timer]
OnCalendar=*:00:00
Unit=pgsql_backup.service

[Install]
WantedBy=timers.target
